# config/supervisor/avahi.conf
[program:avahi]
priority=10
startretries=20
directory=/tmp
command=/usr/local/bin/avahi.sh
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=KILL
stopwaitsecs=1
