# config/supervisor/resolver.conf
[program:resolver]
priority=10
directory=/opt/resolver
command=/usr/bin/npm start
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=KILL
stopwaitsecs=1
